# debian/postinst
#!/bin/sh
# postinst script for netbeans-php-support
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

setxdebug() {
    XDI=$1
        grep -q -e 'xdebug.remote_enable' $XDI || echo 'xdebug.remote_enable = 1' >> $XDI
        grep -q -e 'xdebug.remote_handler' $XDI || echo 'xdebug.remote_handler=dbgp' >> $XDI
        grep -q -e 'xdebug.remote_mode' $XDI || echo 'xdebug.remote_mode=req' >> $XDI
        grep -q -e 'xdebug.remote_host' $XDI || echo 'xdebug.remote_host=127.0.0.1' >> $XDI
        grep -q -e 'xdebug.remote_port' $XDI || echo 'xdebug.remote_port=9000' >> $XDI
        grep -q -e 'xdebug.file_link_format' $XDI || echo 'xdebug.file_link_format = "netbeans://%f?line=%l"' >> $XDI
        grep -q -e 'xdebug.profiler_enable' $XDI || echo 'xdebug.profiler_enable = 0' >> $XDI
#        grep -q -e '' $XDI || echo '' >> $XDI        
}

reloadwebserver() {
        if [ -f /etc/init.d/nginx ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d nginx reload 3>/dev/null || true
            else
                /etc/init.d/nginx reload 3>/dev/null || true
            fi
        fi

        if [ -f /etc/init.d/lighttpd ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d lighttpd reload 3>/dev/null || true
            else
                /etc/init.d/lighttpd reload 3>/dev/null || true
            fi
        fi

        if [ -f /etc/init.d/apache2 ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 reload 3>/dev/null || true
            else
                /etc/init.d/apache2 reload 3>/dev/null || true
            fi
        fi
}


case "$1" in
    configure)

        XDI="/etc/php5/mods-available/xdebug.ini"
	if [ -f $XDI ]; then
	    setxdebug $XDI

            if [ -d "/etc/php5/apache2/" ]; then
	        mv  /etc/php5/apache2/php.ini /etc/php5/apache2/php.ini.orig
    	        cp  /usr/share/php5/php.ini-development /etc/php5/apache2/php.ini
            fi    
            
	    if [ -d "/etc/php5/cli/" ]; then
    	        mv  /etc/php5/cli/php.ini /etc/php5/cli/php.ini.orig
        	cp  /usr/share/php5/php.ini-development /etc/php5/cli/php.ini
            fi
        fi

	XDI="/etc/php/5.6/mods-available/xdebug.ini"
        if [ -f $XDI ]; then
    	    setxdebug $XDI
	fi

        XDI="/etc/php/7.0/mods-available/xdebug.ini"
	if [ -f $XDI ]; then
		setxdebug $XDI
        fi

	XDI="/etc/php/7.1/mods-available/xdebug.ini"
        if [ -f $XDI ]; then
		setxdebug $XDI
        fi

	reloadwebserver;

	if [ -f /usr/bin/xdg-desktop-menu ]; then
	    xdg-desktop-menu install --novendor /usr/share/applications/netbeans.desktop > /dev/null || true
	fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

if [ -f /usr/bin/xdg-desktop-menu ]; then
    xdg-desktop-menu install --novendor /usr/share/applications/netbeans.desktop > /dev/null || true
fi


exit 0
